# systemd/canary-compare.service
[Unit]
Description=pi4gpio direct/pi4gpio canary comparison (MIGRATION_PLAN.md section 6)
After=network.target pi4gpio.service sensor-tiered-client.service
Wants=pi4gpio.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi4gpio/scripts
# 本番venv(lgpio/spidev/smbus2/pyserial/RPi.bme280導入済み)のPythonを使う。
# rpi_sensorsは二重モード対応版(~/rpi-sensor-lib-canary)を使い、本番の
# site-packages(二重モード化前の0.1.0)には一切触れない。
ExecStart=/home/pi/sensor-tiered-store/.venv/bin/python3 canary_compare.py --interval 30 --output /home/pi/pi4gpio/scripts/canary_log.csv
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
